# ansible/roles/lead/templates/pi-swarm-orchestrator.service.j2
[Unit]
Description=Pi Swarm Orchestrator
After=network.target

[Service]
User=swarm
WorkingDirectory={{ pi_swarm_app_dir }}
EnvironmentFile=/etc/pi-swarm/gitea.env
Environment=GITEA_API_BASE_URL={{ gitea_api_base_url }}
Environment=SWARM_METRICS_DB=/var/lib/pi-swarm/metrics.db
Environment=SWARM_CONFIG_PATH={{ pi_swarm_app_dir }}/config/swarm_config.yaml
ExecStart=/usr/bin/python3 -m uvicorn src.orchestrator.app:app --host 0.0.0.0 --port 8080
Restart=always
MemoryMax=512M

[Install]
WantedBy=multi-user.target
